# init/lifeguard.conf
#! upstart
start on startup
stop on shutdown

respawn

exec start-stop-daemon --iosched real-time:1 --start --exec /home/clearwater/rjp/l3dash/run.py
